# systemd/aiolos.service
[Unit]
Description=Aiolos — agnostic fan/airflow orchestrator
Documentation=https://github.com/ktsaou/aiolos
After=network.target

[Service]
Type=simple
ExecStart=/opt/aiolos/bin/aiolos
# aiolos and the anemoi require root: /dev/ipmi0 is root-only and NVML fan control needs root.
User=root

# KillMode is the systemd DEFAULT (control-group): on stop, SIGTERM reaches aiolos AND every module
# at once. Each module catches it and restores its own device (BMC auto / firmware) before exiting
# — modules are self-sufficient and do not depend on aiolos orchestrating their shutdown. The unit
# therefore stays simple and makes no assumption about which modules are configured.

# Belt-and-suspenders fail-safe for a HARD kill (SIGKILL / crash / OOM), where modules are killed
# before they can self-restore: aiolos reads its OWN registry and runs each configured module's
# `restore` one-shot. Routing through aiolos keeps this unit agnostic of the module set (no module
# names here). `-` ignores its exit status; aiolos bounds each module's restore internally so a
# wedged BMC/NVML cannot hang the stop. Idempotent.
ExecStopPost=-/opt/aiolos/bin/aiolos restore

# aiolos logs to stderr (tracing). Capture both streams in the journal — do NOT send to /dev/null.
StandardOutput=journal
StandardError=journal
# Dedicated journal namespace: view with `journalctl --namespace=aiolos [-f]`.
LogNamespace=aiolos
Environment=RUST_LOG=info

Restart=on-failure
RestartSec=5s

# On stop, aiolos catches SIGTERM, tells every module to restore its device (BMC auto / firmware
# fan control), waits for them, then exits. Give it room before systemd escalates to SIGKILL.
KillSignal=SIGTERM
# Headroom for the full stop sequence: aiolos graceful shutdown — bounded by max(5s, the largest
# per-module `timeout` + ~2s restore grace), so ~7s with default timeouts (SOW-0013) — THEN the
# ExecStopPost `aiolos restore` one-shots (bounded ~10s). 30s covers the default config with room.
# IMPORTANT: if you raise any module's `timeout=` in aiolos.conf well above the 5s default, raise
# this too (roughly that timeout + ~12s), or systemd may SIGKILL aiolos mid-restore.
TimeoutStopSec=30s

# NOTE: do NOT add MemoryDenyWriteExecute — NVML/driver userspace may map W+X memory (same
# constraint as the C nvfd unit). Other aggressive sandboxing can also block /dev/ipmi0 + NVML.

[Install]
WantedBy=multi-user.target
